# void/services/thinkfan/run
#!/bin/sh

[ -r ./conf ] && . ./conf

exec thinkfan -nq ${OPTS:-} 2>&1
