# artix/services/ananicy-cpp/finish
#!/bin/sh

exec /sbin/sysctl -qe kernel.sched_autogroup_enabled=1
